# rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Starts server to serve the app
# ==============================================================================
# shellcheck disable=SC2164
cd /opt/ip-backend

exec /opt/ip-backend/ip-backend
